# packages/oscam/unit/oscam.service
[Unit]
Description=Open Source CAM
After=syslog.target
ConditionPathExists=/etc/oscam/oscam.conf

[Service]
ExecStartPre=/usr/bin/oscam-pre.sh
ExecStart=/usr/bin/oscam.sh

[Install]
WantedBy=multi-user.target
