# packages/vdr/unit/vdr.target
[Unit]
Description=VDR (Video Disk Recorder)
Requires=graphical.target
Conflicts=rescue.service rescue.target
After=graphical.target
AllowIsolate=yes

[Install]
Alias=default.target
