# packages/vdr/unit/vdr.service
[Unit]
Description=Video Disk Recorder
Wants=xorg.service
After=xorg.service eventlircd.service

[Service]
Environment=DISPLAY=:0.0
ExecStart=/usr/bin/runvdr
Restart=always

[Install]
WantedBy=vdr.target
