# packages/vdr/unit/xmltv2vdr.service
[Unit]
Description=xmltv2vdr
After=vdr.service vdr-backend.service

[Service]
ExecStart=/usr/bin/xmltv2vdr.sh

[Install]
WantedBy=multi-user.target
