# packages/setup/unit/setup.service
[Unit]
Description=Setup
After=systemd-user-sessions.service plymouth-quit-wait.service plymouth-lite.service
Before=getty.target

[Service]
Type=oneshot
ExecStart=/sbin/setup init
StandardInput=tty-force
StandardOutput=inherit
StandardError=inherit

[Install]
WantedBy=setup.target
