# packages/enigma2pc/unit/epg.service
[Unit]
Description=Update EPG
After=oscam.service

[Service]
Type=oneshot
ExecStart=/usr/bin/update-epg
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
