# packages/backend-service/unit/backend.target
[Unit]
Description=TV streaming server
Requires=network.target
Conflicts=rescue.service rescue.target
After=network.target oscam.target
AllowIsolate=yes

[Install]
Alias=default.target
