# packages/automountd/unit/automountd.service
[Unit]
Description=udev automounter daemon
Conflicts=installator.target
After=systemd-udev-settle.service

[Service]
ExecStart=/usr/bin/automountd
ExecStart=/usr/bin/after_automountd
Restart=always

[Install]
WantedBy=multi-user.target
