# config/platforms/arm/imx6/packages/xbmc/unit/xbmc.target
[Unit]
Description=XBMC Media Center
Requires=multi-user.target
Conflicts=rescue.service rescue.target
After=multi-user.target
AllowIsolate=yes

[Install]
Alias=default.target
