# packages/vdradmin/unit/vdradmin.service
[Unit]
Description=A web based user interface to the VDR

[Service]
ExecStart=/usr/bin/vdradmind 

[Install]
WantedBy=multi-user.target
